# images/base/dind/s6-overlay/s6-rc.d/dockerd/run
#!/command/execlineb -P
# Run dockerd with proper logging
fdmove -c 2 1
/usr/local/bin/dockerd